# debian/wlan-poller.upstart
description "wlan-poller"

start on runlevel [23]
stop on shutdown

script
exec 2>>/var/log/poller.log
exec 1>&2
wlan-poller
end script

respawn
