# debian/wlan-poller.upstart
description "wlan-poller"

start on started hostapd
stop on stopped hostapd

script
    exec 2>>/var/log/poller.log
    exec 1>&2
    wlan-poller
end script

respawn
